# com.zaza.order-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zaza.order-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/zifcrypto/Desktop/zaza/scripts/order-sync-launcher.sh</string>
    </array>

    <!-- DST-aware dual schedule for 9:31 AM ET weekdays -->
    <!-- EDT (Mar-Nov): 13:31 UTC = 9:31 AM ET (runs), 14:31 UTC = 10:31 AM ET (exits) -->
    <!-- EST (Nov-Mar): 13:31 UTC = 8:31 AM ET (exits), 14:31 UTC = 9:31 AM ET (runs) -->
    <!-- Python time-gate in __main__.py exits early if not within 9:26-9:36 AM ET -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday 13:31 UTC -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Monday 14:31 UTC -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Tuesday 13:31 UTC -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Tuesday 14:31 UTC -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Wednesday 13:31 UTC -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Wednesday 14:31 UTC -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Thursday 13:31 UTC -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Thursday 14:31 UTC -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Friday 13:31 UTC -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Friday 14:31 UTC -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/zifcrypto/.zaza/logs/order_sync_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zifcrypto/.zaza/logs/order_sync_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/zifcrypto/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
