# com.zaza.order-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zaza.order-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/zifcrypto/.zaza/scripts/order-sync-launcher.sh</string>
    </array>

    <!-- DST-aware dual schedule for 9:31 AM ET weekdays -->
    <!-- NOTE: StartCalendarInterval uses LOCAL time (SGT +08), not UTC -->
    <!-- EDT (Mar-Nov): 9:31 AM ET = 13:31 UTC = 21:31 SGT (runs) -->
    <!-- EST (Nov-Mar): 9:31 AM ET = 14:31 UTC = 22:31 SGT (runs) -->
    <!-- Python time-gate in __main__.py exits early if not within 9:26-9:36 AM ET -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday 21:31 SGT (EDT: 9:31 AM ET) -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Monday 22:31 SGT (EST: 9:31 AM ET) -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Tuesday 21:31 SGT (EDT: 9:31 AM ET) -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Tuesday 22:31 SGT (EST: 9:31 AM ET) -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Wednesday 21:31 SGT (EDT: 9:31 AM ET) -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Wednesday 22:31 SGT (EST: 9:31 AM ET) -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Thursday 21:31 SGT (EDT: 9:31 AM ET) -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Thursday 22:31 SGT (EST: 9:31 AM ET) -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Friday 21:31 SGT (EDT: 9:31 AM ET) -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
        <!-- Friday 22:31 SGT (EST: 9:31 AM ET) -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>31</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/zifcrypto/.zaza/logs/order_sync_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zifcrypto/.zaza/logs/order_sync_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/zifcrypto/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
